# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# PHP 内置服务器 + router.php 白名单过滤
[program:php]
command=php -S 0.0.0.0:8080 -t /var/www/html /var/www/html/router.php
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/php.log
stderr_logfile=/var/log/supervisor/php-error.log

# 注入守护进程（worker）
[program:worker]
command=php /var/www/html/service/worker.php
autostart=true
autorestart=true
priority=20
startsecs=3
startretries=10
stdout_logfile=/var/log/supervisor/worker.log
stderr_logfile=/var/log/supervisor/worker-error.log
